# ops/launchd/com.cexapidocs.sync.overnight.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  SETUP: launchd plists do not support environment variables.
  Before loading, replace ALL occurrences of /Users/YOURUSER/Projects/cex-api-docs
  with your actual repo path. Then:
    cp ops/launchd/com.cexapidocs.sync.overnight.plist ~/Library/LaunchAgents/
    launchctl load ~/Library/LaunchAgents/com.cexapidocs.sync.overnight.plist
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.cexapidocs.sync.overnight</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>-lc</string>
    <string>cd /Users/YOURUSER/Projects/cex-api-docs &amp;&amp; scripts/run_sync_preset.sh overnight-safe ./cex-docs &gt;&gt; logs/launchd-sync.log 2&gt;&amp;1</string>
  </array>

  <key>RunAtLoad</key>
  <false/>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>2</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>

  <key>WorkingDirectory</key>
  <string>/Users/YOURUSER/Projects/cex-api-docs</string>

  <key>StandardOutPath</key>
  <string>/Users/YOURUSER/Projects/cex-api-docs/logs/launchd-sync.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/YOURUSER/Projects/cex-api-docs/logs/launchd-sync.err.log</string>

  <key>ProcessType</key>
  <string>Background</string>
</dict>
</plist>
